# windows/t1053.005-schtasks-userwritable-path/rule.yml
title: Scheduled Task Created via Schtasks Executing From User-Writable Path
id: 7f3a9c41-2b8d-4e06-9a17-c5d2e8f0b364
status: test
description: >
  Detects schtasks.exe creating a scheduled task whose action executes from a
  user-writable directory (AppData, Temp, Public). Common persistence pattern
  after initial access; legitimate software rarely schedules tasks from these
  paths without a signed installer parent.
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: 'Cullen Flynn'
date: 2026/06/12
logsource:
  product: windows
  category: process_creation
detection:
  selection_img:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/create'
  selection_path:
    CommandLine|contains:
      - '\AppData\'
      - '\Temp\'
      - '\Users\Public\'
      - '%APPDATA%'
      - '%TEMP%'
  condition: selection_img and selection_path
falsepositives:
  - Software updaters (Chrome, Teams, Zoom) creating tasks under AppData — tune
    by ParentImage signature status or known-good task name allowlist
level: medium
tags:
  - attack.persistence
  - attack.t1053.005
